# 60-videowall.rules

KERNEL=="sdb1", SYMLINK+="videowall-target", ACTION=="add", RUN+="/home/perry/test.sh"
KERNEL=="sdb1", ACTION=="remove", RUN+="/home/perry/test.sh"
